# rc.d/chaindogs
#!/bin/ksh

daemon="/usr/local/bin/python3"
daemon_flags="/home/rsroka/chain-of-dogs/web.py"
daemon_user="rsroka"
rc_bg=YES

. /etc/rc.d/rc.subr
rc_cmd $1
